# crowcoin/contrib/init/crowcoind.conf
description "Crowcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CROWCOIND_BIN="/usr/bin/crowcoind"
env CROWCOIND_USER="crowcoin"
env CROWCOIND_GROUP="crowcoin"
env CROWCOIND_PIDDIR="/var/run/crowcoind"
# upstart can't handle variables constructed with other variables
env CROWCOIND_PIDFILE="/var/run/crowcoind/crowcoind.pid"
env CROWCOIND_CONFIGFILE="/etc/crowcoin/crowcoin.conf"
env CROWCOIND_DATADIR="/var/lib/crowcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # crowcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CROWCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run crowcoind."
        echo "The setting must appear in $CROWCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Crowcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CROWCOIND_PIDDIR"
    chmod 0755 "$CROWCOIND_PIDDIR"
    chown $CROWCOIND_USER:$CROWCOIND_GROUP "$CROWCOIND_PIDDIR"
    chown $CROWCOIND_USER:$CROWCOIND_GROUP "$CROWCOIND_CONFIGFILE"
    chmod 0660 "$CROWCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CROWCOIND_PIDFILE" \
    --chuid $CROWCOIND_USER:$CROWCOIND_GROUP \
    --exec "$CROWCOIND_BIN" \
    -- \
    -pid="$CROWCOIND_PIDFILE" \
    -conf="$CROWCOIND_CONFIGFILE" \
    -datadir="$CROWCOIND_DATADIR" \
    -disablewallet \
    -daemon

